# systemd/kaoz.service
[Unit]
Description=Kaoz IRC notifier bot
After=network.target

[Service]
ExecStart=/usr/bin/kaoz

[Install]
WantedBy=multi-user.target
